# test-monitor.service
[Unit]
Description=Monitor for process "test"
After=network.target

[Service]
ExecStart=/usr/local/bin/test_monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target
